# overlay-squashfs-x86/etc/init.d/autorun
#!/sbin/openrc-run

depend()
{
	after legacyrescuecd tigervnc mdadm
}

start()
{
	ebegin "Starting LegacyRescueCd autorun scripts"
	export LANG=en_US.utf8
	echo "Starting autorun at $(date +%Y%m%d-%H:%M:%S)" > /var/log/autorun.log
	/usr/sbin/autorun
}
